# systemd/gbrain-dream.service
[Unit]
Description=gbrain nightly dream cycle

[Service]
Type=oneshot
ExecStart=/data/fleet-memory/systemd/gbrain-dream.sh
TimeoutStartSec=600
