# systemd/gbrain-dream.timer
[Unit]
Description=Run gbrain dream cycle nightly at 3am

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
